# contrib/debathena-transform-lighttpd/debian/debathena-transform-lighttpd.postinst
#!/bin/sh
set -e

#DEBHELPER#
/usr/sbin/invoke-rc.d lighttpd stop
/usr/sbin/invoke-rc.d lighttpd start
