# contrib/debathena-transform-rng-tools/debian/debathena-transform-rng-tools.postinst
#!/bin/sh

#DEBHELPER#
/usr/sbin/invoke-rc.d rng-tools stop
/usr/sbin/invoke-rc.d rng-tools start
