# debian/elbe-daemon.service
[Unit]
Description=ELBE Daemon

[Service]
Type=simple
EnvironmentFile=/etc/default/elbe-daemon
ExecStart=/usr/bin/elbe $DAEMON_ARGS
