# debian/python3-elbe-bootup-check.service
[Unit]
Description=ELBE Bootup check

[Service]
Type=simple
ExecStart=/usr/bin/elbe bootup-check

[Install]
WantedBy=multi-user.target
